# systemd/homebridge.service
[Unit]
Description=Homebridge

[Service]
ExecStart=/home/jamie/homebridge/node_modules/.bin/homebridge
Restart=on-failure
